# dynamic-dns.timer
[Unit]
Description=Check for new DNS records or changes in IP and fix them in DigitalOcean

[Timer]
Persistent=true
OnBootSec=1m
